# wildcloud.conf
#  Wildcloud

description	"Wildcloud"

start on runlevel [2345]
stop on runlevel [!2345]

exec /root/wildcloud
